# containers/cockpit.container
[Container]
Image=quay.io/cockpit/ws:latest
Name=cockpit-ws
AutoUpdate=registry
RunArgs=--privileged --pid=host --network=host -v /:/host
